# systemd/qbittorrent-nox.service
[Unit]
Description=qBittorrent Service
After=network.target

[Service]
Type=forking
User=qbt
Group=qbt
<<<<<<< HEAD
UMask=007
=======
UMask=000
>>>>>>> 887aed5 (Add dkms conf for ixgbe, foundry support, qbt-nox)
ExecStart=/usr/bin/qbittorrent-nox -d --webui-port=12345
Restart=on-failure

[Install]
WantedBy=multi-user.target
